# .claude/skills/add-local-browser-tool/files/com.nanoclaw.browser.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nanoclaw.browser</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Applications/Google Chrome.app/Contents/MacOS/Google Chrome</string>
        <string>--remote-debugging-port=9222</string>
        <string>--remote-debugging-address=127.0.0.1</string>
        <string>--remote-allow-origins=*</string>
        <string>--user-data-dir={HOME}/.nanoclaw-browser-profile</string>
        <string>--no-first-run</string>
        <string>--no-default-browser-check</string>
        <string>--disable-features=ChromeWhatsNewUI</string>
        <!--
            We deliberately do NOT set --disable-blink-features=AutomationControlled
            here, even though it'd suppress navigator.webdriver = true. Chrome
            classifies it as an unsupported flag and shows a yellow infobar
            ("Stability and security will suffer") that's both visible to the
            operator and a fingerprintable surface itself. The helper's
            evaluateOnNewDocument shim shadows navigator.webdriver before any
            page script runs, which achieves the same goal without the infobar.
        -->
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/nanoclaw-browser.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/nanoclaw-browser.err.log</string>
</dict>
</plist>
